# liquidsoap.service
[Unit]
Description=Liquidsoap
After=network.target
Documentation=http://liquidsoap.fm/

[Service]
Type=forking
User=liquidsoap
PIDFile=/run/liquidsoap-radio.pid
WorkingDirectory=/etc/liquidsoap/
ExecStart=/usr/bin/liquidsoap /etc/liquidsoap/radio.liq -d
Restart=always
RestartSec=15
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
Alias=liquidsoap.service
