# liquidsoap.service
[Unit]
Description=Liquidsoap
After=network.target
Documentation=http://liquidsoap.fm/

[Service]
Type=forking
User=liquidsoap
PIDFile=/etc/liquidsoap/liquidsoap.pid
WorkingDirectory=/etc/liquidsoap/
ExecStartPre=touch /etc/liquidsoap/liquidsoap.pid
ExecStart=/usr/bin/liquidsoap /etc/liquidsoap/radio.liq
ExecStop=cat /etc/liquidsoap/liquidsoap.pid | xargs kill
Restart=always

[Install]
WantedBy=multi-user.target
Alias=liquidsoap.service
